# apps/raspi/adaptlight.service
[Unit]
Description=AdaptLight Voice-Controlled Lighting
After=network.target dev-ttyAMA0.device
Wants=dev-ttyAMA0.device

[Service]
Type=simple
User=root
WorkingDirectory=/home/lamp/adaptlight
ExecStart=/bin/bash /home/lamp/adaptlight/scripts/run_raspi.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
