# apps/raspi/adaptlight.service
[Unit]
Description=AdaptLight Voice-Controlled Lighting
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/lamp/adaptlight
EnvironmentFile=/home/lamp/adaptlight/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/lamp/adaptlight/venv/bin/python -m raspi.main --verbose
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
